# deploy/costco-api.service
[Unit]
Description=코스트코 수익계산 FastAPI (Phase2)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/costco-app
EnvironmentFile=/opt/costco-app/.env
# 포트 8610 — 8501(Streamlit)·8600(costco-license)과 겹치지 않게. nginx가 /api/* 라우팅(Phase4).
ExecStart=/opt/costco-app/venv/bin/uvicorn api.main:app \
    --host 127.0.0.1 \
    --port 8610 \
    --workers 2

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=costco-api

# 보안 강화
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
